# scs/puppet/scs/templates/nginx/supervisor.conf.erb
[program:nginx]
command = /usr/sbin/nginx -c /scs/etc/nginx/nginx.conf
stdout_logfile = /scs/var/log/supervisord/%(program_name)s-stdout.log
stderr_logfile = /scs/var/log/supervisord/%(program_name)s-stderr.log
stopwaitsecs = 30
