# scs/puppet/scs/templates/php-fpm/supervisor.conf.erb
[program:php-fpm]
command = /usr/sbin/php5-fpm -y /scs/etc/php-fpm/php-fpm.ini
stdout_logfile = /scs/var/log/supervisord/%(program_name)s-stdout.log
stderr_logfile = /scs/var/log/supervisord/%(program_name)s-stderr.log
stopwaitsecs = 30
